# crontab
SHELL=/bin/bash
0 16 * * * make cron -f ./epoch-archive/Makefile | tee -a ./cron.log
# Mainnet epoch happens at approx 15:30 UTC

#* * * * * make cron -f ./epoch-archive/Makefile | tee -a ./cron.log
