# crontab
SHELL=/bin/bash
0 22 * * * make cron -f ./epoch-archive/Makefile | tee -a ./cron.log

#* * * * * make cron -f ./epoch-archive/Makefile | tee -a ./cron.log
